# deploy/crowd-canvas.service
[Unit]
Description=Crowd Canvas
After=network.target

[Service]
Type=simple
# Use a user that owns /opt/crowd-canvas (see README). www-data is a safe default.
User=www-data
Group=www-data
WorkingDirectory=/opt/crowd-canvas
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
# basic hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
